# share/systemd/user/gsel.socket
[Unit]
Description=gsel server

[Socket]
ListenStream=@%t/gsel.sock

[Install]
WantedBy=default.target

